# deploy/systemd/temasotoyolyardim-gunicorn.service.example
[Unit]
Description=Temas Oto Yol Yardim — Django REST API (Gunicorn WSGI)
Documentation=https://docs.djangoproject.com/en/stable/howto/deployment/wsgi/gunicorn/
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data

WorkingDirectory=/home/ubuntu/Tow_Truck_Temas/backend
EnvironmentFile=-/home/ubuntu/Tow_Truck_Temas/backend/.env
Environment=DJANGO_SETTINGS_MODULE=config.settings.prod
Environment=APP_ROOT=/home/ubuntu/Tow_Truck_Temas

ExecStart=/home/ubuntu/Tow_Truck_Temas/venv/bin/gunicorn \
    --config /home/ubuntu/Tow_Truck_Temas/deploy/gunicorn/temasotoyolyardim.py \
    config.wsgi:application

RuntimeDirectory=temasotoyolyardim
RuntimeDirectoryMode=0750

Restart=on-failure
RestartSec=5

# Büyük dosya yüklemeleri (nginx client_max_body_size ile uyumlu)
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
