# deploy/systemd/temasotoyolyardim-astro-rebuild.service.example
[Unit]
Description=Temas Oto Yol Yardim — Astro rebuild webhook (admin kaydı → npm run build:prod)
Documentation=file:///home/ubuntu/Tow_Truck_Temas/PROJECT.md
After=network.target temasotoyolyardim-gunicorn.service
Wants=temasotoyolyardim-gunicorn.service

[Service]
Type=simple
# npm/astro build ubuntu ile çalışmalı; frontend/ root’a ait olmamalı (chown ubuntu:ubuntu frontend)
User=ubuntu
Group=ubuntu
SupplementaryGroups=www-data

WorkingDirectory=/home/ubuntu/Tow_Truck_Temas/frontend
EnvironmentFile=-/home/ubuntu/Tow_Truck_Temas/backend/.env
Environment=WEBHOOK_PORT=9876
Environment=ASTRO_BUILD_SCRIPT=build:prod
Environment=BUILD_API_URL=http://127.0.0.1:8765
# BUILD_API_URL → webhook PUBLIC_BUILD_API_URL (SSG fetch). PUBLIC_API_URL .env.production'da kalır (form/API).
# PROD_ASTRO_REBUILD_WEBHOOK_SECRET backend/.env dosyasından okunur (EnvironmentFile)

ExecStart=/usr/bin/node /home/ubuntu/Tow_Truck_Temas/frontend/scripts/webhook-rebuild.mjs

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
