# deploy/systemd/temasotoyolyardim-gunicorn.service.example
[Unit]
Description=Temas Oto Yol Yardim — Django REST API (Gunicorn WSGI)
Documentation=https://docs.djangoproject.com/en/stable/howto/deployment/wsgi/gunicorn/
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data

WorkingDirectory=/home/ubuntu/Tow_Truck_Temas/backend
# systemd .env okur (root); Gunicorn worker'ları da python-dotenv ile okuyabilsin:
#   sudo chmod 640 /home/ubuntu/Tow_Truck_Temas/backend/.env
#   sudo chown root:www-data /home/ubuntu/Tow_Truck_Temas/backend/.env
# Şifrede * # $ varsa tırnaksız yazın. Satır sonu # yorumu YAZMAYIN (systemd değere karıştırır).
EnvironmentFile=-/home/ubuntu/Tow_Truck_Temas/backend/.env
Environment=DJANGO_SETTINGS_MODULE=config.settings.prod
Environment=APP_ROOT=/home/ubuntu/Tow_Truck_Temas

# Groq arka plan çevirisi backend/logs/ altına yazar (www-data yazabilmeli).
ExecStartPre=/bin/mkdir -p /home/ubuntu/Tow_Truck_Temas/backend/logs
ExecStartPre=/bin/chown www-data:www-data /home/ubuntu/Tow_Truck_Temas/backend/logs
ExecStartPre=/bin/chmod 750 /home/ubuntu/Tow_Truck_Temas/backend/logs

ExecStart=/home/ubuntu/Tow_Truck_Temas/venv/bin/gunicorn \
    --config /home/ubuntu/Tow_Truck_Temas/deploy/gunicorn/temasotoyolyardim.py \
    config.wsgi:application

RuntimeDirectory=temasotoyolyardim
RuntimeDirectoryMode=0750

Restart=on-failure
RestartSec=5

# Büyük dosya yüklemeleri (nginx client_max_body_size ile uyumlu)
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
